# deploy/linux/xmrig-watchdog.timer
[Unit]
Description=Run the XMRig busy-watchdog every 30s

[Timer]
OnBootSec=2min
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target
