# deploy/linux/xmrig.service.d/10-forge.conf
# sigma-forge only: it is a BUILD node on a 15 W mobile i7-8650U that throttles under
# sustained load. CPUWeight=1 + Nice=19 (in the base unit) already make the miner yield
# instantly to builds; this hard ceiling additionally leaves ~1 logical core cool for build
# responsiveness and caps sustained package heat. Remove/raise if you want forge to mine harder
# when fully idle. (studio, a desktop part, does NOT get this drop-in.)
[Service]
CPUQuota=300%
